# pi-media-hub.service.template
[Unit]
Description=Pi Media Hub
After=network.target

[Service]
Type=simple
User=REPLACE_USER
WorkingDirectory=REPLACE_DIR/pi-media-hub
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
